# package/custom/sing-box/files/sing-box.init
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
PROG=/usr/bin/sing-box
CONF=/etc/sing-box/config.json

start_service() {
	procd_open_instance
	procd_set_param command $PROG run -c $CONF -D /var/lib/sing-box
	procd_set_param user sing-box
	procd_set_param file $CONF
	procd_set_param respawn
	procd_close_instance
}
